# deploy/social-surveyor@.service
[Unit]
Description=social-surveyor (%i)
Documentation=https://github.com/apurvam/social-surveyor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=social-surveyor
Group=social-surveyor
WorkingDirectory=/opt/social-surveyor
EnvironmentFile=/etc/social-surveyor/%i.env

# SQLite DB and any per-project state live outside the repo so rsync/git
# deploys don't touch them.
Environment="SOCIAL_SURVEYOR_DATA_DIR=/var/lib/social-surveyor/%i"

# SSM Parameter Store fallback for secrets.resolve_secret. When set,
# the resolver tries this SSM prefix if a name isn't in the process
# env. Belt-and-suspenders alongside the EnvironmentFile path: the
# bootstrap's social-surveyor-load-env helper pre-fills the env file
# from SSM, so steady-state hits the env path; SSM is the fallback
# for names the operator forgot to add to SSM / the env file, or
# secrets rotated without re-seeding.
Environment="SOCIAL_SURVEYOR_SSM_PREFIX=/social-surveyor/%i"

# ExecStartPre with '+' prefix runs as root (needed for chown).
ExecStartPre=+/bin/mkdir -p /var/lib/social-surveyor/%i
ExecStartPre=+/bin/chown -R social-surveyor:social-surveyor /var/lib/social-surveyor/%i

ExecStart=/usr/local/bin/uv run --directory /opt/social-surveyor social-surveyor run --project %i

Restart=on-failure
RestartSec=30s

# journald captures stdout/stderr; view via `journalctl -u social-surveyor@%i`
StandardOutput=journal
StandardError=journal
SyslogIdentifier=social-surveyor-%i

# Resource ceilings for a single t4g.micro running one project.
MemoryMax=256M
TasksMax=32

[Install]
WantedBy=multi-user.target
